# roles/openclaw/templates/xvfb.service.j2
[Unit]
Description=Xvfb Virtual Framebuffer
After=network.target

[Service]
Type=simple
User={{ openclaw_user }}
ExecStart=/usr/bin/Xvfb :99 -screen 0 1920x1080x24
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
